# anacrontab
# $HOME/.anacron/anacrontab: configuration file for anacron

# See anacron(8) and anacrontab(5) for details.

SHELL=/bin/sh
HOME=/home/david
LOGNAME=user
MAILTO="user"

# These replace cron's entries
# days	delay	job-name	command
1	5	cron.daily	run-parts --report $HOME/.anacron/cron.daily
7	10	cron.weekly	run-parts --report $HOME/.anacron/cron.weekly
@monthly	15	cron.monthly	run-parts --report $HOME/.anacron/cron.monthly
