# cron.weekly/anacron_user_log_cleanup
#!/bin/bash

LOG_FILE="anacron.log"
TEMP_FILE="anacron.log.tmp"
WEEK_AGO=$(date -d "7 days ago" +%Y-%m-%d)

# 1. Identify and keep at least the last 5 entries for major categories
# Regardless of age, ensuring you don't lose recent history for critical jobs.
(
	grep "cron.yearly" "$LOG_FILE" | tail -n 5
	grep "cron.monthly" "$LOG_FILE" | tail -n 5
	grep "cron.weekly" "$LOG_FILE" | tail -n 5
    # Keep last 5 for daily jobs
    #~ grep "cron.daily" "$LOG_FILE" | tail -n 5
    # Keep last 5 for hourly jobs
    #~ grep "cron.hourly" "$LOG_FILE" | tail -n 5

    # 2. Add all entries from the last week that aren't already captured
    awk -v limit="$WEEK_AGO" '$1 >= limit' "$LOG_FILE"
) > "$TEMP_FILE"

# Replace the old log with the cleaned version
mv "$TEMP_FILE" "$LOG_FILE"
for file in ls spool/anacron; do
	if cat anacrontab | grep -v "^#" | awk "{print $3}"| grep $file -q;then
		continue
	else
		rm spool/anacron/$file
	fi
done
exit 0